# packaging/systemd/vfio-sensor-bridge-agent.service
[Unit]
Description=VFIO Sensor Bridge Linux guest agent
After=dev-virtio\x2dports-org.vfio_sensor_bridge.0.device

[Service]
Type=simple
ExecStart=/usr/sbin/agent-linux run --config /etc/vfio-sensor-bridge/agent.toml
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
